# docker/supervisor/conf.d/filebeat.conf
[program:filebeat]
command=/usr/share/filebeat/filebeat -c /etc/filebeat/filebeat.yml
startsecs=0
autostart=true
autorestart=true
