# docker/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
startsecs=0
autostart=true
autorestart=true
